# docker/php/crontab
* * * * * /usr/local/bin/php /var/www/html/artisan schedule:run > /var/www/html/storage/logs/cron.last.log
